# psiphon-fleet.service
[Unit]
Description=Psiphon Multi-Instance Fleet
Documentation=https://github.com/rezasmind/x-ui-pro
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/psiphon-fleet
ExecStart=/usr/bin/docker compose -f docker-compose-psiphon.yml up -d
ExecStop=/usr/bin/docker compose -f docker-compose-psiphon.yml down
ExecReload=/usr/bin/docker compose -f docker-compose-psiphon.yml restart
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=psiphon-fleet

[Install]
WantedBy=multi-user.target
